# launchagents/com.untionglim.memory-guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.untionglim.memory-guardian</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/untionglim/.local/bin/memory-guardian</string>
    </array>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/untionglim/.local/log/memory-guardian.err</string>
</dict>
</plist>
